# systemd-garrett-units/samba.service
[Unit]
Description=Samba

[Service]
Type=forking
ExecStart=/etc/rc.d/samba start
ExecStop=/etc/rc.d/samba stop

[Install]
WantedBy=multi-user.target
